# pi-image/files/beacon.service
[Unit]
Description=Beacon Tally
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Restart=on-failure
RestartSec=3
WorkingDirectory=/opt/Beacon-tally/app
Environment=NODE_ENV=production
ExecStart=/usr/local/bin/node /opt/Beacon-tally/app/dist/supervisor.js

[Install]
WantedBy=multi-user.target
