# pi-image/files/beacon.service
[Unit]
Description=Beacon Tally
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Restart=on-failure
RestartSec=3
WorkingDirectory=/opt/Beacon-tally/app
Environment=NODE_ENV=production
ExecStart=/opt/Beacon-tally/app/node_modules/.bin/tsx /opt/Beacon-tally/app/supervisor.ts

[Install]
WantedBy=multi-user.target
